# m/data/config/org.mongodb.default.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.mongodb.default</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/mongod</string>
    <string>run</string>
    <string>--dbpath</string>
    <string>/data/db/current/default</string>
    <string>--logpath</string>
    <string>/data/db/current/default/mongod.log</string>
    <string>--logappend</string>
    <string>--bind_ip</string>
    <string>127.0.0.1</string>
    <string>--auth</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>UserName</key>
  <string>brandon</string>
  <key>WorkingDirectory</key>
  <string>/data/db/current/default</string>
  <key>StandardErrorPath</key>
  <string>/data/db/current/default/mongod.log</string>
  <key>StandardOutPath</key>
  <string>/data/db/current/default/mongod.log</string>
  <key>SoftResourceLimits</key>
  <integer>64000</integer>
</dict>
</plist>
